# deploy/systemd/bms-province-discovery-full.service
[Unit]
Description=BMS Province Discovery FULL sweep (safety net for incremental discovery)
After=network.target

[Service]
Type=oneshot
User=bms
Group=bms
WorkingDirectory=/opt/bms/app
EnvironmentFile=/opt/bms/app/.env
ExecStart=/opt/bms/venv/bin/python /opt/bms/app/scripts/Sebastian_Province_Discovery.py --worker --ingest --full
TimeoutStartSec=1800
Nice=10
StandardOutput=journal
StandardError=journal
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
