# deploy/systemd/bms-deadman.service
[Unit]
Description=BMS Dead-Man Switch (token + discovery health)
After=network.target

[Service]
Type=oneshot
User=bms
Group=bms
WorkingDirectory=/opt/bms/app
EnvironmentFile=/opt/bms/app/.env
ExecStart=/opt/bms/venv/bin/python /opt/bms/app/scripts/health_deadman.py
Nice=10
StandardOutput=journal
StandardError=journal
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
