# deploy/systemd/bms-province-discovery-full.timer
[Unit]
Description=Daily BMS Province FULL sweep (00:30 UTC = 07:30 Thai) — reconcile incremental

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true
Unit=bms-province-discovery-full.service

[Install]
WantedBy=timers.target
